# .config/systemd/user/sync-to-git.timer
[Unit]
Description=sync notes to git

[Timer]
OnBootSec=3min
OnCalendar=hourly
Unit=sync-to-git.target

[Install]
WantedBy=timers.target

